# crane/upstart/docker-apt-cacher-ng.conf
description "apt-cacher-ng"
author "Roman Zubatyuk"
start on filesystem and started docker
stop on runlevel [!2345]
respawn
script
     cd %INSTALLDIR% 
     exec crane run apt-cacher-ng
end script
pre-stop script
    /usr/bin/docker stop apt-cacher-ng
end script
